# deploy/com.leo.cortex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.leo.cortex</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/leoyuan/Desktop/leoworkspace/cortex/.venv/bin/python</string>
    <string>-m</string>
    <string>cortex.main</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/leoyuan/Desktop/leoworkspace/cortex</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/leoyuan/Desktop/leoworkspace/cortex/data/stdout.log</string>
  <key>StandardErrorPath</key><string>/Users/leoyuan/Desktop/leoworkspace/cortex/data/stderr.log</string>
</dict>
</plist>
